# src/roles/daily_report/templates/daily-report.service.j2
[Unit]
Description=Daily system status report

[Service]
Type=oneshot
ExecStart={{ daily_report_script_path }}
