# src/roles/docker_prune/templates/docker-prune.timer.j2
[Unit]
Description=Schedule Docker System Prune
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnCalendar={{ docker_prune_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
