# src/roles/docker_prune/templates/docker-prune.service.j2
[Unit]
Description=Docker System Prune
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart={{ docker_prune_command }}
