# src/roles/backup_restic/templates/rpi-backup.timer.j2
[Unit]
Description=Schedule Restic Backup
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnCalendar={{ backup_restic_timer_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
